# supervisord.conf
[supervisord]
nodaemon=true

[program:sshd]
command=/usr/sbin/sshd -D
autostart=true

[program:jira]
command=/usr/local/atlassian/jira/bin/start-jira.sh -fg
user=daemon
